# mise/tasks/vm/files/idle-shutdown.timer
[Unit]
Description=Periodic idle check for automatic shutdown

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min
Unit=idle-shutdown.service

[Install]
WantedBy=timers.target
