# systemd/flow-gazette@.service
[Unit]
Description=Flow Broker (%i)
Documentation=https://github.com/estuary/flow
After=flow-etcd.service
Requires=flow-etcd.service

[Service]
Type=simple

# Load instance-specific environment
# Instance format: <cluster>-<port> (e.g., foobar-8080)
EnvironmentFile=%h/flow-local/env/common.env
EnvironmentFile=%h/flow-local/env/gazette-%i.env

ExecStart=%h/flow-local/run-in-repo go run go.gazette.dev/core/cmd/gazette serve

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
