# systemd/flow-plane@.target
[Unit]
Description=Flow Data Plane (%i)
Documentation=https://github.com/estuary/flow
StopWhenUnneeded=yes

[Install]
WantedBy=default.target
