# local/systemd/flow-dekaf@.service
[Unit]
Description=Flow Dekaf (%i)
Documentation=https://github.com/estuary/flow
# After=/Requires=flow-dekaf-kafka@<stack>.service is added as a drop-in by the
# data-plane task, which knows this data plane's stack (dekaf-kafka is keyed by
# stack name, not by the data-plane instance name here).

[Service]
Type=simple
TimeoutStartSec=600

EnvironmentFile=%h/flow-local/env/dekaf-%i.env

ExecStart=sh -c 'exec ${CARGO_TARGET_DIR}/debug/dekaf \
    --local \
    --advertise-host 127.0.0.1 \
    --kafka-port ${DEKAF_KAFKA_PORT} \
    --schema-registry-port ${DEKAF_REGISTRY_PORT} \
    --metrics-port ${DEKAF_METRICS_PORT} \
    --default-broker-urls tcp://localhost:${DEKAF_UPSTREAM_PORT} \
    --upstream-auth none \
    --encryption-secret ${DEKAF_ENCRYPTION_SECRET} \
    --data-plane-fqdn ${DATA_PLANE_FQDN} \
    --data-plane-access-key ${DATA_PLANE_ACCESS_KEY} \
    --task-refresh-interval 10s \
    --spec-ttl 10s \
    --idle-session-timeout 60s'

Restart=on-failure
RestartSec=5s
