# local/systemd/flow-gazette@.service
[Unit]
Description=Flow Broker (%i)
Documentation=https://github.com/estuary/flow
# After=flow-etcd@<stack>.service is added as a drop-in by the gazette task,
# which knows this instance's stack (the etcd instance is keyed by stack name,
# not by the <dp>-<port> instance name here).

[Service]
Type=simple
TimeoutStartSec=600

# Single consolidated env file (common build/runtime vars inlined at generation).
# Instance format: <cluster>-<port> (e.g., foobar-8080)
EnvironmentFile=%h/flow-local/env/gazette-%i.env

ExecStartPre=sh -c 'cd ${FLOW_ROOT} && ./mise/tasks/build/gazette'
# Two CORS allow-origins (gazette's BROKER_ALLOW_ORIGIN env var takes only one):
# the global external UI, plus the optional per-stack UI at base+5. Both vars come
# from the EnvironmentFile above.
ExecStart=sh -c 'exec ${GOBIN}/gazette serve --broker.allow-origin ${FLOW_DASHBOARD_ORIGIN} --broker.allow-origin http://localhost:${FLOW_PORT_DASHBOARD}'

Restart=on-failure
RestartSec=5s
