# local/systemd/flow-bigtable.service
[Unit]
Description=Flow BigTable Emulator (local development)
Documentation=https://github.com/estuary/flow

[Service]
Type=simple

ExecStart=docker run --rm --name flow-bigtable \
  -p 8086:8086 \
  google/cloud-sdk:latest \
  gcloud beta emulators bigtable start --host-port=0.0.0.0:8086

ExecStop=-/bin/sh -c 'docker stop flow-bigtable >/dev/null 2>&1 || true'

Restart=on-failure
RestartSec=5s
